# worker-sw/worker.service
[Unit]
Description=wgmesh Worker Service
After=network.target

[Service]
Type=simple
Environment=WORKER_API_KEY=your_api_key_here
ExecStart=/usr/local/bin/wgmesh-worker
Restart=on-failure

[Install]
WantedBy=multi-user.target
